# SD3/supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:nginx]  
command=nginx
  
[program:sd3]  
command=php /apps/bin/start_swoole_server.php start -d
